# provisioning/deploy/files/prowler.timer
[Unit]
Description=Run prowler once a day on weekdays

[Timer]
OnCalendar=Mon..Fri 12:00
Persistent=false

[Install]
WantedBy=timers.target